# systemd/send-observation.timer
[Unit]
Description=Send observed sensor data periodically

[Timer]
# Every 5 minutes
OnCalendar=*:0/5

[Install]
WantedBy=timers.target
